# supervisord.conf
[supervisord]
nodaemon=true
[program:cron]
command=cron -f
[program:syslog]
command=service rsyslog start
[program:updatettnode]
command=/usr/node/ttnode_task.sh update
